# src/rc.local
#!/bin/bash

source /etc/drone/env.sh

$MAIN_ENTRY

exit 0