# src/rc.local
#!/bin/bash

source /etc/dronevery/env.sh

$MAIN_ENTRY

exit 0